# src/build/docker/cloudrouter-bird-centos/assets/supervisord.conf
[supervisord]
nodaemon=true

[program:bird]
command=/usr/sbin/bird -f -c /etc/bird/bird.conf
