# src/build/docker/cloudrouter-bird-fedora/assets/supervisord.conf
[supervisord]
nodaemon=true

[program:bird]
command=/usr/sbin/bird -d -c /etc/bird/bird.conf
